# deploy/target/etc/systemd/system/evkkapi.service
[Unit]
Description=EVKK API
After=syslog.target network.target postgresql.service redis-server.service

[Service]
Type=simple

User=evkk
Group=evkk
EnvironmentFile=/etc/environment
EnvironmentFile=/opt/evkk/conf/common.env
EnvironmentFile=/opt/evkk/conf/api.env
WorkingDirectory=/opt/evkk/service/api
ExecStart=/usr/bin/java -jar api.jar
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
