# deploy/target/etc/systemd/system/evkkclusterfinder.service
[Unit]
Description=EVKK CLUSTERFINDER
After=syslog.target network.target postgresql.service redis-server.service

[Service]
Type=simple

User=evkk
Group=evkk
EnvironmentFile=/etc/environment
EnvironmentFile=/opt/evkk/conf/common.env
EnvironmentFile=/opt/evkk/conf/clusterfinder.env
WorkingDirectory=/opt/evkk/service/clusterfinder
ExecStart=/usr/bin/java -jar clusterfinder.jar
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
