# deploy/target/etc/systemd/system/evkkdaemon.service
[Unit]
Description=EVKK DAEMON
After=syslog.target network.target postgresql.service redis-server.service

[Service]
Type=simple

User=evkk
Group=evkk
EnvironmentFile=/etc/environment
EnvironmentFile=/opt/evkk/conf/common.env
EnvironmentFile=/opt/evkk/conf/daemon.env
WorkingDirectory=/opt/evkk/service/daemon
ExecStart=/usr/bin/java -jar daemon.jar
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
